# sender/supervisor/postfix.conf
[program:postfix]
command=
autorestart=true